# netlocation/netlocation@.service
[Unit]
Description=Network geo-location service %m %H:%i

Requires=etcd.service
Requires=docker.service

After=etcd.service
After=docker.service

[Service]
# From: https://www.digitalocean.com/community/tutorials/how-to-create-flexible-services-for-a-coreos-cluster-with-fleet-unit-files
# Change killmode from "control-group" to "none" to let Docker remove
# work correctly.
KillMode=none

EnvironmentFile=/etc/environment
EnvironmentFile=/home/core/share/adNimbusEnvironment

# Pulling the container the first time can take more than the default 90secs. Set to 0 for no timeout
TimeoutStartSec=0

ExecStartPre=-/usr/bin/docker rm -f ${netLocationDockerTag}_%i
ExecStartPre=/home/core/share/adnimbus_registry/startAdNimbusRegistry.sh startPre ${netLocationService}
ExecStartPre=/usr/bin/docker ps -a

ExecStart=/home/core/share/netlocation/startNetlocation.sh start %i

ExecStartPost=/home/core/share/netlocation/startNetlocation.sh registerNetLocation %i

ExecStop=/home/core/share/netlocation/startNetlocation.sh cleanup %i

#Restart=always

