# monitor/monitor@.service
[Unit]
Description=Work with consul to monitor services across the clusters. Precursor to using Salt, Ngios, Ganglia and/or Sensu.

[Service]
# Allow access to environment variables.
EnvironmentFile=/etc/environment
EnvironmentFile=/home/core/ad-nimbus/adNimbusEnvironment
EnvironmentFile=/home/core/ad-nimbus/monitor/monitorEnvironment

User=core

ExecStart=/home/core/ad-nimbus/monitor/monitor.sh start

ExecStop=/home/core/ad-nimbus/monitor/monitor.sh stop

